# data/rc.local
#!/bin/bash

/home/jxminer/setup/scripts/10-dhcp-client.sh

exit 0
